# systemd/electricitystatistics.service
[Unit]
Description=recording electricity statistics

[Service]
#Type=simple
ExecStart=/usr/bin/python3 /home/ubuntu/energy_statistic.py
